# backend/animereloaded-mal-auth.service
[Unit]
Description=AnimeReloaded MyAnimeList Auth Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/animereloaded-mal-auth
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/animereloaded-mal-auth.env
ExecStart=/usr/bin/python3 /opt/animereloaded-mal-auth/backend/animereloaded_mal_auth_server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
